# projects/api/rootfs/etc/services.d/php_fpm/run
#!/usr/bin/with-contenv sh
set -e;

# Start PHP-FPM
exec /usr/sbin/php-fpm7 -R --nodaemonize
